# scripts/com.hxy.clauded.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key><string>com.hxy.clauded</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{REPO}}/.venv/bin/clauded</string>
    </array>
    <key>WorkingDirectory</key><string>{{REPO}}</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key><string>{{HOME}}</string>
        <key>PATH</key><string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>RunAtLoad</key><true/>
    <key>KeepAlive</key><true/>
    <key>ThrottleInterval</key><integer>30</integer>
    <key>StandardOutPath</key><string>{{HOME}}/Library/Logs/clauded/out.log</string>
    <!-- #168 dropped: ``StandardErrorPath`` previously pointed at err.log, which
         mirrored every Python log line (Python's stderr handler in bot.py +
         launchd's redirect both wrote the same content). RotatingFileHandler
         at ~/Library/Logs/clauded/clauded.log is now the canonical log. -->
    <!-- #232 round 3: ProcessType removed entirely (default=Standard).
         Round 1 removed Background, Round 2 set Interactive — both still
         triggered "because inefficient" kills. Round 3 removes the key
         entirely and adds LowPriorityBackgroundIO=false. Empirically
         verified: 60+ min stable with 0 kills (was 15-30s between kills
         with Interactive). The exact launchd-internal mechanism is not
         fully documented; this combination works on macOS 15.x. -->
    <key>LowPriorityBackgroundIO</key><false/>
</dict>
</plist>
