# scripts/com.hxy.clauded.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key><string>com.hxy.clauded</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{REPO}}/.venv/bin/clauded</string>
    </array>
    <key>WorkingDirectory</key><string>{{REPO}}</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key><string>{{HOME}}</string>
        <key>PATH</key><string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>RunAtLoad</key><true/>
    <key>KeepAlive</key><true/>
    <key>ThrottleInterval</key><integer>30</integer>
    <key>StandardOutPath</key><string>{{HOME}}/Library/Logs/clauded/out.log</string>
    <!-- #168 dropped: ``StandardErrorPath`` previously pointed at err.log, which
         mirrored every Python log line (Python's stderr handler in bot.py +
         launchd's redirect both wrote the same content). RotatingFileHandler
         at ~/Library/Logs/clauded/clauded.log is now the canonical log. -->
    <key>ProcessType</key><string>Interactive</string>
    <!-- #232 follow-up (round 2): empty ProcessType still trips
         launchd's default "light resource limits + throttling" for
         daemon-spawned services (per launchd.plist(5)). The
         "because inefficient" SIGKILL kept happening even with the
         key removed. `Interactive` declares this is a foreground
         interactive service — launchd skips the resource judgement
         entirely. See docs/prd/v1.18-fix-launchagent-process-type.md
         for the original investigation. -->
    <!-- #232 removed: ``ProcessType=Background`` told launchd this is a
         low-priority occasional task, so it SIGKILL'd us every 15-25 min
         with the "because inefficient" reason. clauded is a long-running
         interactive service (Discord gateway + SDK long-poll). Falling
         back to the implicit default (Standard) gives us no launchd-side
         resource judgement. See docs/prd/v1.18-fix-launchagent-process-type.md. -->
</dict>
</plist>
